# docker/php/crontab
PATH=/usr/local/bin:/usr/bin:/bin:/usr/local/sbin:/usr/sbin:/sbin
* * * * * root /script.sh
# An empty line is required at the end of this file for a valid cron file.